# scripts/systemd/acme.socket
[Socket]
Accept=no
ListenStream=45.83.105.255:8080
Service=acme-proxy.service

[Install]
WantedBy=sockets.target