# systemd/acme.socket
[Socket]
Accept=no
ListenStream=8080
Service=acme-proxy.service

[Install]
WantedBy=sockets.target

